# mysql_base/supervisord.conf
[supervisord]
nodaemon=true

[program:mysql]
command=/usr/sbin/mysqld

[program:tomcat]
command=/usr/local/mysql/default/bin/mysqld_safe --initialize --console
user=root